# supervisord.conf
[supervisord]
nodaemon=true

[program:supervisor-touchstone]
command=python3 -c "while True: print('')" >> /dev/null

